# zettai-mamorukun.service
[Unit]
Description=ぜったいまもるくん - サイバー攻撃防御デーモン
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/zettai-mamorukun --config /etc/zettai-mamorukun/config.toml
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
